# deploy/systemd/colorado-traffic-tracker.service
[Unit]
Description=Colorado Traffic Tracker Docker Compose stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/colorado-traffic-tracker
Environment=APP_ENV_FILE=.env.cloud
ExecStart=/usr/bin/docker compose --env-file .env.cloud up -d --remove-orphans
ExecStop=/usr/bin/docker compose --env-file .env.cloud down
RemainAfterExit=yes
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
